# rded/service_files/rded.service
[Unit]
Description=Phosphor RDE Daemon

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/rded

[Install]
WantedBy=multi-user.target
